# packages/pivirt/pivirt.initd
#!/sbin/openrc-run

name="pivirt"
description="pivirt hypervisor daemon"
command="/usr/bin/myhypervisor"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
    need net
    after firewall
}
